# infrastructure/edge/gitops/edge-gitops-sync.timer
# /etc/systemd/system/edge-gitops-sync.timer
# =============================================================================
# Polls the homelab GitHub repo every 5 minutes and deploys changed configs.
# Persistent: catches up if the machine was off at the scheduled time.
#
# To trigger immediately after a push:
#   sudo systemctl start edge-gitops-sync.service
# =============================================================================
[Unit]
Description=GitOps config sync timer for edge node (every 5 min)

[Timer]
OnBootSec=60
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
