# infrastructure/vault/seal-vault/vault.service
# Vault Seal Server — systemd Unit  
# Deploy to: /etc/systemd/system/vault.service on 192.168.20.28
[Unit]
Description=HashiCorp Vault — Seal Server (Transit Auto-Unseal Provider)
Documentation=https://developer.hashicorp.com/vault/docs
Requires=network-online.target
After=network-online.target
# Must start BEFORE production Vault nodes
Before=

[Service]
User=vault
Group=vault
ExecStart=/usr/bin/vault server -config=/etc/vault.d/vault.hcl
ExecReload=/bin/kill --signal HUP $MAINPID
KillMode=process
KillSignal=SIGINT
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
LimitNOFILE=65536
LimitMEMLOCK=infinity

# Security hardening
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_IPC_LOCK
AmbientCapabilities=CAP_IPC_LOCK

[Install]
WantedBy=multi-user.target
